# mail-filter/spampd/files/init-r1
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

piddir="/var/run"
pidfile="${piddir:-/var/run}/spampd.pid"
required_dirs="${piddir:-/var/run}"
command="/usr/sbin/spampd.pl"
command_args="--pid=${pidfile} ${SPAMPD_OPTS:-}"

depend() {
	use net
}
