# net-analyzer/ntop/files/ntop-initd-r1
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="ntop - a network traffic analyzer"
pidfile="/var/run/ntop.pid"
command="/usr/bin/ntop"
command_args="-L ${NTOP_OPTS}"
start_stop_daemon_args="-p ${pidfile} -m -b"

depend() {
	need net

	if [[ -e /etc/conf.d/tmpfs ]]; then
		config /etc/conf.d/tmpfs
		source /etc/conf.d/tmpfs
	fi

	local x="RC_BEFORE_${RC_SVCNAME}"
	[[ -n ${!x} ]] && before ${!x}
	local x="RC_NEED_${RC_SVCNAME}"
	[[ -n ${!x} ]] && need ${!x}
	local x="RC_USE_${RC_SVCNAME}"
	[[ -n ${!x} ]] && use ${!x}
	local x="RC_AFTER_${RC_SVCNAME}"
	[[ -n ${!x} ]] && after ${!x}
}

start_pre() {
	if [ ! -e /var/lib/ntop/ntop_pw.db ] ; then
		eerror "You need to set a password first by running"
		eerror "ntop --set-admin-password"
		return 1
	fi

	export LC_ALL=C # apparently doesn't work with some locales (#191576 and #205382)
	export GDFONTPATH=/usr/share/fonts/corefonts # apparently the only font that works (#231705)
}

start() {

	ebegin "Starting ntop"
	#start-stop-daemon --start --exec /usr/bin/ntop --ionice 3 --nice 2 -- -d -L ${NTOP_OPTS}
	start-stop-daemon --start --exec ${command} ${start_stop_daemon_args} \
		--ionice 3 --nice 2 -- ${command_args}
	eend $?
	if (( ! $? )) && type -pf cpulimit >/dev/null 2>&1; then
		sleep 0.1
		local count=0
		einfo "Waiting for 'ntop' to initialise ..."
		while true; do
			if ! [[ -e "${pidfile}" ]]; then
				einfo "Waiting for PID file '${pidfile}' ..."
				sleep 0.1
				if [[ -e "${pidfile}" ]]; then
					continue
				fi
				(( count++ ))
				if (( count > 100 )); then
					eerror "... timed out waiting for PID, aborting"
					return 1
				fi
			elif ! ps $( cat "${pidfile}" ) >/dev/null 2>&1; then
				einfo "Waiting for process '$( cat "${pidfile}" )' ..."
				sleep 0.1
				if ps $( cat "${pidfile}" ) >/dev/null 2>&1; then
					einfo "... found"
					count=0
					break
				fi
				(( count++ ))
				if (( count > 100 )); then
					eerror "... timed out waiting for process, aborting"
					return 1
				fi
			else
				break
			fi
		done
		ebegin "Starting cpulimit for ntop($( cat "${pidfile}" ))"
		start-stop-daemon --start --exec $( type -pf cpulimit ) -- -p $( cat "${pidfile}" ) -l 20 -b -z
		eend $?
	fi
}
