# sys-fs/lvm2/files/device-mapper.rc-2.02.97-r1
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-fs/lvm2/files/device-mapper.rc-2.02.95-r2,v 1.2 2013/04/09 11:00:26 ssuominen Exp $

depend() {
	before dmeventd checkfs fsck
	after modules
}

start() {
	local f="/lib/rcscripts/addons/dm-start.sh"

	if [ -r "${f}" ]; then
		( . "${f}" )
	fi
}

