# sys-fs/lvm2/files/lvmetad.initd-2.02.105-r2
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id: 2f9ab0c3d9f7676283ec15771a21fbd5d1cb89c0 $

pidfile="/var/run/lvmetad.pid"
command="/sbin/lvmetad"
command_args="${LVMETAD_OPTS:=-p ${pidfile}}"
start_stop_daemon_args="--pidfile ${pidfile}"

start_pre() {
	checkpath -q -d /var/run/lvm
}
