# sys-fs/lvm2/files/lvmetad.initd-2.02.116-r3
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

pidfile="/var/run/lvmetad.pid"
command="/sbin/lvmetad"
command_args="${LVMETAD_OPTS:=-p ${pidfile}}"
start_stop_daemon_args="--pidfile ${pidfile}"

depend() {
	:
}

start_pre()
{
	checkpath -quiet --directory /var/run/lvm || return 1
}
