# sys-apps/kmod/files/kmod-static-nodes-r1
#!/sbin/openrc-run
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id: c3eeb2aeb1f15c1a72e95f69785a4963266d5ffc $

description="Create list of required static device nodes for the current kernel"

depend() {
	after dev-mount
	before tmpfiles.dev dev
	keyword -lxc -systemd-nspawn
}

start() {
	ebegin "Creating list of required static device nodes for the current kernel"
	checkpath -q -d /var/run/tmpfiles.d
	kmod static-nodes --format=tmpfiles --output=/var/run/tmpfiles.d/kmod.conf
	eend $?
}
