# net-analyzer/munin/files/munin-graph.socket
[Unit]
Description=Munin FastCGI Graph Socket

[Socket]
ListenStream=/var/run/munin/fcgi-graph.sock

[Install]
WantedBy=sockets.target
