# sys-apps/busybox/files/klogd.initd
#!/sbin/openrc-run
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id: 99f5f7b05467c66866f30f25dc64a30803128938 $

command="/bin/busybox klogd"
command_args="${KLOGD_OPTS}"
pidfile="/var/run/klogd.pid"

depend() {
	need clock hostname
}
