# sys-apps/busybox/files/watchdog.initd
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id: aa1d33eb7c06cfec232270cba8354fad61e4b1d5 $

command="/bin/busybox watchdog"
command_args="${WATCHDOG_OPTS}"
pidfile="/var/run/watchdog.pid"
