# sys-power/apcupsd/files/apcupsd.init.4
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id: 2a69687db89da28955a9567dd4c6084f68c42908 $

INSTANCE="${SVCNAME#*.}"
if [ -z "${INSTANCE}" ] || [ "${SVCNAME}" = "apcupsd" ]; then
	INSTANCE="apcupsd"
fi

piddir="/var/run/apcupsd"

depend() {
	use net
}

start_pre() {
	local lockdir="$( grep -E '^\s*LOCKFILE\s+' "/etc/apcupsd/${INSTANCE}.conf" | sed -r 's/\s*LOCKFILE\s+//' )"
        checkpath -d -m 0775 -o root:uucp ${lockdir} 
        checkpath -d -m 0775 -o root:uucp ${piddir} 
}

start() {
	rm -f /etc/apcupsd/powerfail

	export SERVICE="${SVCNAME}"

	ebegin "Starting APC UPS daemon"
	start-stop-daemon \
		--start --pidfile "${piddir}/${SVCNAME}.pid" \
		--exec /sbin/apcupsd -- \
		-f "/etc/apcupsd/${INSTANCE}.conf" \
		-P "${piddir}/${SVCNAME}.pid"
	eend $?
}

stop() {
	ebegin "Stopping APC UPS daemon"
	start-stop-daemon \
		--stop --pidfile "${piddir}/${SVCNAME}.pid" \
            	--retry TERM/5/TERM/5 \
		--exec /sbin/apcupsd
	eend $?
}

