# net-analyzer/rrdtool/files/rrdcached.init
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id: 364695ee5a9b9fe6c66bc45cd3fdcc4ea7f37b95 $

command="/usr/bin/rrdcached"
command_args="${RRCACHE_ARGS} -p /var/run/rrdcached.pid"
start_stop_daemon_args="--quiet"
description="RRDtool data caching daemon"

depend() {
	need localmount net
	use logger
}
